# docker/supervisord.base.conf
; =====================================================================================
; Glenn Runtime Base — Supervisor configuration
; =====================================================================================
; Minimal base config. The daemon (once installed at /opt/agent/) is started here; all
; project-specific services (postgres + anything the spec adds: redis, minio, mailhog,
; mongodb, dev servers, etc.) are added at bootstrap time by writing fragments into
; /data/.glenn/supervisor.d/ and calling `supervisorctl reread && supervisorctl update`.
;
; Per Runtime Spec V2, only Postgres has its binary pre-baked in the base image; all
; other services depend on their spec's per-service `install` script having run during
; bootstrap before supervisord tries to exec them. The daemon's bootstrap pipeline
; orders InstallStage before the supervisor fragments are written, so by the time
; `supervisorctl update` fires, the binaries exist.
; =====================================================================================

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/tmp/supervisord.pid
user=agent

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; --------------------------------------------------------------------------------------
; The agent daemon. The daemon binary is NOT baked into the base image — it is
; downloaded at boot from R2 by /usr/local/bin/bootstrap-daemon.sh, which reads
; DAEMON_BUNDLE_URL + DAEMON_BUNDLE_SHA256 (stamped in by RuntimeProvisionerJob from
; the /api/daemon-versions/resolve endpoint), sha256-verifies the tarball, extracts to
; /opt/agent/, and exec's `node /opt/agent/daemon.js`.
;
; Other env contract (also stamped by main API at Machine create time):
;   GLENN_RUNTIME_TOKEN, MAIN_API_URL, RUNTIME_ID, DAEMON_VERSION.
;
; startsecs is 30 because cold start has to do curl + sha256 + tar-extract before node
; even starts; 5s wasn't enough headroom. startretries=3 because if R2 is unreachable
; we want to fail fast rather than crash-loop forever.
; --------------------------------------------------------------------------------------
[program:agent]
command=/usr/local/bin/bootstrap-daemon.sh
directory=/data
autostart=true
autorestart=true
startretries=3
startsecs=30
stopwaitsecs=10
stopsignal=TERM
stdout_logfile=/var/log/supervisor/agent.out.log
stderr_logfile=/var/log/supervisor/agent.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
environment=NODE_ENV="production",MISE_DATA_DIR="/data/mise"

; --------------------------------------------------------------------------------------
; Include drop-in fragments. Bootstrap writes per-project services (postgres etc.)
; into this directory, then runs `supervisorctl reread && update` — supervisor picks
; them up without restart.
; --------------------------------------------------------------------------------------
[include]
files = /data/.glenn/supervisor.d/*.conf
