# .agents/skills/agent-dev-services/templates/dotnet-api.conf
[program:dotnet-api]
command=/data/project/repo/.claude/skills/agent-dev-services/scripts/run-dotnet-api.sh
directory=/data/project/repo/packages/dotnet-api
user=agent
autostart=true
autorestart=true
startsecs=20
startretries=5
stopwaitsecs=15
stopsignal=TERM
stdout_logfile=/var/log/supervisor/dotnet-api.out.log
stderr_logfile=/var/log/supervisor/dotnet-api.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
environment=HOME="/home/agent",DOTNET_ROOT="/home/agent/.dotnet",PATH="/home/agent/.dotnet:/home/agent/.dotnet/tools:/usr/share/dotnet:/usr/bin:/bin",DATABASE_URL="Host=localhost;Port=5432;Database=app;Username=postgres;Password=postgres",ASPNETCORE_ENVIRONMENT="Development"
