# .claude/skills/agent-dev-services/templates/backoffice-web.conf
[program:backoffice-web]
command=/data/project/repo/.claude/skills/agent-dev-services/scripts/run-backoffice-web.sh
directory=/data/project/repo/packages/backoffice-web
user=agent
autostart=true
autorestart=true
startsecs=5
startretries=5
stopwaitsecs=10
stopsignal=TERM
stdout_logfile=/var/log/supervisor/backoffice-web.out.log
stderr_logfile=/var/log/supervisor/backoffice-web.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
environment=HOME="/home/agent",PATH="/usr/bin:/bin",NODE_ENV="development"
